# pi-app/systemd/hotsport-access.service
[Unit]
Description=Hotsport Access Pi-Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
NotifyAccess=main
User=root
Group=root
WorkingDirectory=/opt/hotsport-access/current
Environment=HOTSPORT_ACCESS_CONFIG=/etc/hotsport-access/config.toml
ExecStart=/opt/hotsport-access/venv/bin/python -m app.main
Restart=always
RestartSec=2s

# Daemon sendet alle ~paar Sekunden WATCHDOG=1; bleibt das aus, killt
# systemd den Prozess und startet ihn neu.
WatchdogSec=60s

# Reader hängt bei keyboard/qr im blockierenden read() bis zum nächsten
# Tastendruck. Daemon hat einen 3-s-Force-Exit nach SIGTERM; falls der
# nicht greift, lässt systemd ihn nach 5 s mit SIGKILL beenden statt
# 90 s lang auf einen sauberen Stop zu warten.
TimeoutStopSec=5s
KillSignal=SIGTERM
SendSIGKILL=yes

Nice=-5

# State-Dir muss beschreibbar bleiben, auch wenn / read-only ist
ReadWritePaths=/var/lib/hotsport-access /var/log/hotsport-access

[Install]
WantedBy=multi-user.target
